# ops/systemd/ticketing-compose.service
[Unit]
Description=Ticketing monorepo stack via docker compose
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_ROOT__
ExecStart=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
