# ops/backup/crontab
# Run backups daily at 03:10
10 3 * * * /backup.sh >> /var/log/backup.log 2>&1
